# vector/vector.service
[Unit]
Description=Vector event pipeline
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/ec2-user/zip-external/.env
ExecStart=/usr/bin/vector --config /home/ec2-user/zip-external/vector/vector.yml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target